# bootc/files/bootc-sync-esp.service
[Unit]
Description=Sync bootc bootloader to EFI System Partition
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc-sync-esp
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
